# arbo.service
[Unit]
Description=Arbo Polymarket Trading System
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=arbo
Group=arbo
WorkingDirectory=/opt/arbo
ExecStart=/opt/arbo/.venv/bin/python -m arbo.main --mode paper
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=arbo

# Environment
EnvironmentFile=/opt/arbo/.env

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/arbo/models /opt/arbo/logs

# Memory limit — alert if approaching 4GB instance limit
MemoryMax=3G
MemoryHigh=2G

[Install]
WantedBy=multi-user.target
