# arbo.service
[Unit]
Description=Arbo Polymarket Trading System
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=arbo
Group=arbo
WorkingDirectory=/opt/arbo
ExecStart=/opt/arbo/.venv/bin/python -m arbo --mode rdh
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=arbo

# Environment
EnvironmentFile=/opt/arbo/.env

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/opt/arbo/models /opt/arbo/logs /home/arbo/.cache /home/arbo/.postgresql

# Memory limit — e5-large-v2 model uses ~1.3GB, uvicorn adds ~100MB
MemoryMax=3.5G
MemoryHigh=3G

[Install]
WantedBy=multi-user.target
