# scripts/systemd/arbo-summary-yearly.service
[Unit]
Description=Arbo yearly retrospective
After=network-online.target postgresql.service

[Service]
Type=oneshot
User=arbo
Group=arbo
WorkingDirectory=/opt/arbo
EnvironmentFile=/opt/arbo/.env
ExecStart=/opt/arbo/.venv/bin/python /opt/arbo/scripts/retrospective.py yearly
StandardOutput=journal
StandardError=journal
TimeoutStartSec=1800
