# scripts/systemd/arbo-capture-d-orderbook.service
[Unit]
Description=Arbo — Strategy D orderbook snapshot capture
Documentation=file:/opt/arbo/scripts/capture_d_orderbook.py
After=network.target postgresql.service arbo.service
Wants=postgresql.service
# We don't require main arbo to be running, just prefer after it starts

[Service]
Type=simple
User=arbo
Group=arbo
WorkingDirectory=/opt/arbo
ExecStart=/opt/arbo/.venv/bin/python /opt/arbo/scripts/capture_d_orderbook.py --interval 300 --window-hours 48
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=arbo-capture-d

# Environment (DB credentials + API keys)
EnvironmentFile=/opt/arbo/.env

# Security hardening — same pattern as main arbo.service
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
# Writable paths: logs + potential __pycache__ regen
ReadWritePaths=/opt/arbo/arbo /opt/arbo/logs

# Memory limit — lightweight (network + DB only, no ML model loaded)
MemoryMax=512M
MemoryHigh=400M

[Install]
WantedBy=multi-user.target
